# roles/laptop/templates/iptables_config.service.j2
[Unit]
Description=configure_iptables

[Service]
Type=simple
ExecStart=/opt/net.cdaniels/bin/system-configure-iptables
