# roles/laptop/templates/suspend_prevent_atboot.service.j2
[Unit]
Description=Determine during boot if inhibit lock needs to be set or not
After=dbus.service

[Service]
ExecStart=/opt/net.cdaniels/bin/system-suspendctl

[Install]
WantedBy=multi-user.target
